# templates/vault/vault.service.j2
[Unit]
Description=Vault server
After=network.target

[Service]
ExecStart={{vault_path}}/vault-server {{vault_path}}/vault-config.hcl
StandardOutput=journal

[Install]
WantedBy=multi-user.target

